# deploy/ai-trading.service
[Unit]
Description=AI Trading Bot
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/ai-trading-team
Environment="PATH=/home/ubuntu/.local/bin:/usr/local/bin:/usr/bin:/bin"

# 使用 uv 运行
ExecStart=/home/ubuntu/.local/bin/uv run python main.py

# 自动重启
Restart=always
RestartSec=10

# 日志
StandardOutput=append:/home/ubuntu/ai-trading-team/logs/service.log
StandardError=append:/home/ubuntu/ai-trading-team/logs/service-error.log

# 资源限制
MemoryMax=2G
CPUQuota=80%

[Install]
WantedBy=multi-user.target
